# root/etc/s6-overlay/s6-rc.d/init-config/up
# This file doesn't do anything, it's just the start of the downstream image init process
