# package/docker/root/etc/s6-overlay/s6-rc.d/svc-php-fpm/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ -f /usr/sbin/php-fpm83 ]]; then
    exec /usr/sbin/php-fpm83 -F
fi
