# package/docker/root/etc/s6-overlay/s6-rc.d/init-version-checks/up
/etc/s6-overlay/s6-rc.d/init-version-checks/run
