# root/etc/s6-overlay/s6-rc.d/init-mods/up
# This file doesn't do anything, it's just the start of the mod init process
