# package/docker/root/etc/s6-overlay/s6-rc.d/ci-service-check/up
echo -e " Core         : Completed loading container"
