# root/etc/s6-overlay/s6-rc.d/init-config-end/up
# This file doesn't do anything, it's just the end of the downstream image init process
