# root/etc/s6-overlay/s6-rc.d/init-folders/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# make folders
mkdir -p \
    /config/{keys,php,www} \
    /config/log/{nginx,php} \
    /config/nginx/site-confs \
    /run \
    /var/lib/nginx/tmp/client_body \
    /var/tmp/nginx
