# package/docker/root/etc/s6-overlay/s6-rc.d/init-nginx-end/up
# This file doesn't do anything, it's just the end of the nginx base init process
